# .config/systemd/user/swayidle.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/swayidle -w timeout 1201 'niri msg action power-off-monitors' timeout 1200 'swaylock -f -i ~/Pictures/wallpapers/black-rock-shooter.jpg' before-sleep 'swaylock -f -i ~/Pictures/wallpapers/black-rock-shooter.jpg'
Restart=on-failure
